# Bagel/payload/Library/LaunchAgents/zip.recyclebin.bagel.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>zip.recyclebin.bagel</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/bagel-scan</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>14</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/bagel.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/bagel.stderr.log</string>
</dict>
</plist>
